# config/etc/cron.d/garagesecurity
40 5 * * * /usr/bin/curl -u "username:password" -H "Content-Type: application/json" -XPUT -d '{ "silence": true }' http://localhost/alerts30 8 * * 1-5 /usr/bin/curl -u "username:password" -H "Content-Type: application/json" -XPUT -d '{ "silence": false }' http://localhost/alerts30 15 * * 1-5 /usr/bin/curl -u "username:password" -H "Content-Type: application/json" -XPUT -d '{ "silence": true }' http://localhost/alerts