# config/etc/cron.daily/motion_expire
#!/bin/sh

/usr/bin/find /home/motion -ctime +14 -delete
/usr/bin/find /home/motion -type d -empty -delete
